# rootfs/etc/s6-overlay/s6-rc.d/40-check-ais-notify/up
#!/bin/sh
exec /etc/s6-overlay/scripts/40-check-ais-notify
