# rootfs/etc/s6-overlay/s6-rc.d/80-tropoalert-init/up
#!/bin/sh
exec /etc/s6-overlay/scripts/80-tropoalert-init
